# contrib/init/coinuad.conf
description "CoinUA Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env COINUAD_BIN="/usr/bin/coinuad"
env COINUAD_USER="coinua"
env COINUAD_GROUP="coinua"
env COINUAD_PIDDIR="/var/run/coinuad"
# upstart can't handle variables constructed with other variables
env COINUAD_PIDFILE="/var/run/coinuad/coinuad.pid"
env COINUAD_CONFIGFILE="/etc/coinua/coinua.conf"
env COINUAD_DATADIR="/var/lib/coinuad"

expect fork

respawn
respawn limit 5 120
kill timeout 600

pre-start script
    # this will catch non-existent config files
    # coinuad will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$COINUAD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run coinuad."
        echo "The setting must appear in $COINUAD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"CoinUA Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$COINUAD_PIDDIR"
    chmod 0755 "$COINUAD_PIDDIR"
    chown $COINUAD_USER:$COINUAD_GROUP "$COINUAD_PIDDIR"
    chown $COINUAD_USER:$COINUAD_GROUP "$COINUAD_CONFIGFILE"
    chmod 0660 "$COINUAD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$COINUAD_PIDFILE" \
    --chuid $COINUAD_USER:$COINUAD_GROUP \
    --exec "$COINUAD_BIN" \
    -- \
    -pid="$COINUAD_PIDFILE" \
    -conf="$COINUAD_CONFIGFILE" \
    -datadir="$COINUAD_DATADIR" \
    -disablewallet \
    -daemon

